# deployment/roles/api-server/templates/api-server.service.j2
[Unit]
Description=Phenopod API Server
Documentation=https://github.com/VarmaMSP/cello
After=network.target

[Service]
Type=simple
User=root
ExecStart={{ server.build_dir }}/main
Restart=on-failure

[Install]
WantedBy=multi-user.target
